# packages/app-admin/logrotate/files/logrotate.cron
#! /bin/sh

/usr/bin/logrotate /etc/logrotate.conf
