# packages/net-fs/nfs-utils/files/systemd/nfs-server.service
# Copyright 2010-2011 Wulf C. Krueger <philantrop@exherbo.org>
# Distributed under the terms of the GNU General Public License v2

[Unit]
Description=NFS server
Requires=rpcmountd.service rpcbind.service rpcstatd.service proc-fs-nfsd.mount
After=network.target rpcmountd.service rpcbind.service rpcstatd.service rpcidmapd.service proc-fs-nfsd.mount

[Service]
Type=oneshot
EnvironmentFile=/etc/conf.d/nfs-server.conf
ExecStartPre=/usr/sbin/exportfs -ra
ExecStart=/usr/sbin/rpc.nfsd ${RPC_NFSD_OPTS}
ExecStop=/usr/sbin/exportfs -ua
ExecStop=/usr/sbin/rpc.nfsd 0
RemainAfterExit=yes

[Install]
WantedBy=network.service multi-user.target
