# packages/net-misc/openssh/files/systemd/sshd.service
[Unit]
Description=OpenSSH server daemon
After=syslog.target network.target auditd.service

[Service]
ExecStartPre=/usr/host/bin/sshd-keygen.sh
ExecStart=/usr/host/bin/sshd -D
ExecReload=/usr/host/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
