# packages/net-misc/openssh/files/systemd/sshd-keygen.service
[Unit]
Description=OpenSSH Server Key Generation
ConditionPathExists=|!/etc/ssh/ssh_host_rsa_key
ConditionPathExists=|!/etc/ssh/ssh_host_ed25519_key

[Service]
ExecStart=/usr/host/bin/ssh-keygen -A
Type=oneshot
