# packages/net-nds/rpcbind/files/systemd/rpcbind.service
[Unit]
Description=RPC bind service
Requires=rpcbind.socket

[Service]
Type=forking
ExecStart=/sbin/rpcbind -w

[Install]
Also=rpcbind.socket
