# packages/net-nds/rpcbind/files/systemd/rpcbind.service
[Unit]
Description=RPCbind Service
Requires=rpcbind.socket

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/rpcbind.conf
ExecStart=/usr/bin/rpcbind -w $RPCBIND_ARGS

[Install]
Also=rpcbind.socket
