# packages/sys-libs/glibc/files/systemd/nscd.service
[Unit]
Description=nscd

[Service]
Type=forking
PIDFile=/run/nscd/nscd.pid
ExecStart=/usr/host/bin/nscd
ExecStop=/usr/host/bin/nscd --shutdown

[Install]
WantedBy=multi-user.target
