# packages/sys-apps/xinetd/files/systemd/xinetd.service
[Unit]
Description=Extended Internet service daemon
After=network.target

[Service]
ExecStart=/usr/sbin/xinetd -dontfork

[Install]
WantedBy=multi-user.target
