# packages/sys-fs/mdadm/files/systemd/mdadm.service
[Unit]
Description=Software RAID monitoring and management
After=syslog.target
ConditionPathExists=/etc/mdadm.conf

[Service]
Type=forking
PIDFile=/run/mdadm/mdadm.pid
ExecStart=/sbin/mdadm --monitor --scan -f --pid-file=/run/mdadm/mdadm.pid

[Install]
WantedBy=multi-user.target
